# gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn --worker-class eventlet -w 1 main:app_ -b 0.0.0.0:5000
directory=/app/