# supervisord.conf
[supervisord]
nodaemon=true
childlogdir=/logs

[program:nginx]
command=/usr/sbin/nginx